# gohookr.service
[Unit]
Description=A really simple webhook receiver.
After=network.target

[Service]
ExecStart=/service/alvps/gohookr
